# supervisor/spfw.conf
[program:spfw]
command=/usr/bin/python /opt/spfw/ip.py runserver
autostart=true
autorestart=true
